# example/build/supervisor-confd/screenrec.conf
[program:screenrec]
command=/bin/bash /opt/screenrec.sh
priority=30
user=mainuser
