# example/build/supervisor-confd/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -De
priority=10
