# build/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;

[group:main]
programs=Xvfb,LXDE,x11vnc,noVNC

[program:Xvfb]
command=/usr/bin/Xvfb :0 -screen 0 %resolution%x24 -nolisten tcp -nolisten unix
priority=10
user=%mainuser_name%
environment=USER="%mainuser_name%",HOME="%mainuser_home%"

[program:LXDE]
command=/usr/bin/startlxde
priority=15
user=%mainuser_name%
environment=USER="%mainuser_name%",HOME="%mainuser_home%",DISPLAY=":0"
directory=%mainuser_home%

[program:x11vnc]
command=/usr/bin/x11vnc -display :0 -xkb -forever -shared -repeat -capslock
    -rfbport %vnc_port% %vncpwoption%
priority=20
user=%mainuser_name%
environment=USER="%mainuser_name%",HOME="%mainuser_home%"

[program:noVNC]
command=/usr/bin/websockify --web=/usr/share/novnc %novnc_port%
    127.0.0.1:%vnc_port%
priority=25
user=%mainuser_name%
environment=USER="%mainuser_name%",HOME="%mainuser_home%"

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;

[include]
files = /etc/supervisor/conf.d/*.conf
